# docs/pir-server.service
[Unit]
Description=Combined PIR server (nullifier + witness)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pir-server
ExecStart=/opt/pir-server/pir-server \
    --lwd-url https://zec.rocks:443 \
    --data-dir /opt/pir-server/data \
    --listen 127.0.0.1:8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
