# docs/spend-server.service
[Unit]
Description=Zcash PIR server (nullifier + witness)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/spend-server
ExecStart=/opt/spend-server/spend-server \
    --lwd-url https://zec.rocks:443 \
    --data-dir /opt/spend-server/data \
    --listen 127.0.0.1:8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
